# webhook-listener.service.example
[Unit]
Description=GitHub Webhook Auto-Deploy Listener
After=network.target

[Service]
Type=simple
User=pm2
WorkingDirectory=/home/pm2/github-webhooks
ExecStart=/usr/bin/node /home/pm2/github-webhooks/webhook-server.js
Restart=on-failure
EnvironmentFile=/home/pm2/github-webhooks/.env

[Install]
WantedBy=multi-user.target
