# webhook-listener.service.example
[Unit]
Description=GitHub Webhook Auto-Deploy Listener
After=network.target

[Service]
Type=simple
User=pm2
WorkingDirectory=/var/www/webhook-deploy-listener
ExecStart=/usr/bin/node /var/www/webhook-deploy-listener/webhook-server.js
Restart=on-failure
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
